# config/frpc.supervisor.web.broadcast.conf
[program:frpc_http_broadcast]
command=/root/center_frpc_http/bin/frpc broadcast
stopasgroup=true
killasgroup=true
restart=always
user=root
