# config/frpc.supervisor.web.mobile.conf
[program:mobile_web]
command=/root/center_frpc_http/bin/frpc mobile
killasgroup=true
restart=always
user=root
