# config/frpc.supervisor.web.crawl.conf
[program:frpc_http_crawl]
command=/root/center_frpc_http/bin/frpc crawl
stopasgroup=true
killasgroup=true
restart=always
user=root
