# config/frpc.supervisor.web.argue_system.conf
[program:argue_system_web]
command=/root/center_frpc_http/bin/frpc argue_system
stopasgroup=true
killasgroup=true
restart=always
user=root
