# config/frpc.supervisor.web.neo.conf
[program:frpc_http_neo]
command=/root/center_frpc_http/bin/frpc neo
stopasgroup=true
killasgroup=true
restart=always
user=root
