# config/frpc.supervisor.web.rss.conf
[program:frpc_http_rss]
command=/root/center_frpc_http/bin/frpc rss
stopasgroup=true
killasgroup=true
restart=always
user=root
