# config/frpc.supervisor.web.argue_system.conf
[program:argue_web]
command=/root/center_frpc_http/bin/frpc argue
stopasgroup=true
killasgroup=true
restart=always
user=root
