# config/frpc.supervisor.web.map.conf
[program:map_web]
command=/root/center_frpc_http/bin/frpc map
killasgroup=true
restart=always
user=root
