# config/frpc.supervisor.web.portal.conf
[program:rss_web]
command=/root/center_frpc_http/bin/frpc portal
restart=always
user=root
