# config/frpc.supervisor.web.banner_system.conf
[program:frpc_http_banner]
command=/root/center_frpc_http/bin/frpc banner
stopasgroup=true
killasgroup=true
restart=always
user=root
