# config/frpc.supervisor.web.portal.conf
[program:portal_web]
command=/root/center_frpc_http/bin/frpc portal
killasgroup=true
restart=always
user=root
