# config/frpc.supervisor.web.data_reporter.conf
[program:data_reporter_web]
command=/root/center_frpc_http/bin/frpc data_reporter
restart=always
user=root
