# config/frpc.supervisor.web.neo4j.conf
[program:frpc_http_neo4j]
command=/root/center_frpc_http/bin/frpc neo4j
stopasgroup=true
killasgroup=true
restart=always
user=root
