# config/frpc.supervisor.web.stream.conf
[program:stream_web]
command=/root/center_frpc_http/bin/frpc stream
killasgroup=true
restart=always
user=root
