# config/frpc.supervisor.web.stream.conf
[program:frpc_http_stream]
command=/root/center_frpc_http/bin/frpc stream
stopasgroup=true
killasgroup=true
restart=always
user=root
