# config/frpc.supervisor.web.dialogue.conf
[program:dialogue_web]
command=/root/center_frpc_http/bin/frpc dialogue
killasgroup=true
restart=always
user=root
