# config/frpc.supervisor.web.foodnutrition.conf
[program:foodnutrition_web]
command=/root/center_frpc_http/bin/frpc foodnutrition
killasgroup=true
restart=always
user=root
