# config/frpc.supervisor.web.dialogue.conf
[program:frpc_http_dialogue_web]
command=/root/center_frpc_http/bin/frpc dialogue
stopasgroup=true
killasgroup=true
restart=always
user=root
