# config/frpc.supervisor.web.mobile.conf
[program:frpc_http_mobile]
command=/root/center_frpc_http/bin/frpc mobile
stopasgroup=true
killasgroup=true
restart=always
user=root
