# config/frpc.supervisor-rss_web.conf
[program:rss_web]
command=/root/center_frpc_http/bin/frpc rss.yao-yang.cn
restart=always
user=root
