# config/frpc.supervisor.web.search.conf
[program:frpc_http_search]
command=/root/center_frpc_http/bin/frpc search
stopasgroup=true
killasgroup=true
restart=always
user=root
