# config/frpc.supervisor.web.question.conf
[program:frpc_http_question_web]
command=/root/center_frpc_http/bin/frpc question
stopasgroup=true
killasgroup=true
restart=always
user=root
