# config/frpc.supervisor.web.complain.conf
[program:frpc_http_complain]
command=/root/center_frpc_http/bin/frpc complain
stopasgroup=true
killasgroup=true
restart=always
user=root
