# config/frpc.supervisor.web.rss.conf
[program:rss_web]
command=/root/center_frpc_http/bin/frpc rss
killasgroup=true
restart=always
user=root
