# config/frpc.supervisor.web.map.conf
[program:frpc_http_map]
command=/root/center_frpc_http/bin/frpc map
stopasgroup=true
killasgroup=true
restart=always
user=root
