# config/frpc.supervisor.web.www.conf
[program:portal_web]
command=/root/center_frpc_http/bin/frpc www
stopasgroup=true
killasgroup=true
restart=always
user=root
